# .devcontainer/crontab
# 每 6 小时重启 open-webui 容器
0 */6 * * * /app/restart_container.sh >> /var/log/cron/cron.log 2>&1
